# workspace/crontab/laradock
* * * * * laradock php /var/www/cheap/artisan schedule:run >> /dev/null 2>&1
